# config/com.secondbrain.memory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.secondbrain.memory</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/davidd/second-brain/scripts/memory-service.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/davidd/second-brain</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/davidd/.memory/logs/service-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/davidd/.memory/logs/service-stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>MEMORY_HOME</key>
        <string>/Users/davidd/.memory</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>LowPriorityBackgroundIO</key>
    <true/>
    
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>